# deploy/systemd/novelstudio-healthcheck.timer
[Unit]
Description=Run NovelStudio health checks every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
Persistent=true
Unit=novelstudio-healthcheck.service

[Install]
WantedBy=timers.target
