# deploy/systemd/novelstudio-backup.timer
[Unit]
Description=Run NovelStudio backup every night

[Timer]
OnCalendar=*-*-* 03:17:00
Persistent=true
Unit=novelstudio-backup.service

[Install]
WantedBy=timers.target
