# deploy/systemd/novelstudio-compose.service
[Unit]
Description=NovelStudio Compose Stack
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/root/NovelStudio
Environment=COMPOSE_PROJECT_NAME=novelstudio
ExecStart=/usr/bin/docker compose --env-file /root/NovelStudio/.env.compose up -d --build
ExecStop=/usr/bin/docker compose --env-file /root/NovelStudio/.env.compose down
ExecReload=/usr/bin/docker compose --env-file /root/NovelStudio/.env.compose up -d --build
TimeoutStartSec=0

[Install]
WantedBy=multi-user.target
